# dot_config/systemd/user/rclone@.service.tmpl
[Unit]
Description=Rclone mount %I driectory

[Service]
Type=simple
ExecStart=/usr/bin/rclone mount OneDrive:/%I/  {{ .chezmoi.homeDir }}/%I --cache-dir {{ .chezmoi.homeDir }}/.cache/OneDrive --vfs-cache-mode writes

[Install]
WantedBy=default.target
