# dot_config/systemd/user/wayland-instance@.target
[Unit]
Description=Services which should run after the graphical session is brought up
Requires=basic.target
After=graphical-session.target
